# supervisord.hf.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/var/log/supervisor

; ── Redis ────────────────────────────────────────────────────
[program:redis]
command=redis-server --port 6379 --daemonize no --protected-mode no --loglevel warning
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.log
priority=10

; ── Go API server ────────────────────────────────────────────
[program:api]
command=/usr/local/bin/api
environment=REDIS_ADDR="localhost:6379"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api.log
priority=30
startsecs=3

; ── Go workers (2 replicas via two program sections) ─────────
[program:worker1]
command=/usr/local/bin/worker
environment=REDIS_ADDR="localhost:6379",WORKER_POOL_SIZE="2"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/worker1.log
stderr_logfile=/var/log/supervisor/worker1.log
priority=40

[program:worker2]
command=/usr/local/bin/worker
environment=REDIS_ADDR="localhost:6379",WORKER_POOL_SIZE="2"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/worker2.log
stderr_logfile=/var/log/supervisor/worker2.log
priority=40

; ── Python RAG / FastAPI service ─────────────────────────────
[program:rag]
command=uvicorn main:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app/rag
environment=REDIS_ADDR="localhost:6379",GROQ_API_KEY="%(ENV_GROQ_API_KEY)s",OPENROUTER_API_KEY="%(ENV_OPENROUTER_API_KEY)s"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/rag.log
stderr_logfile=/var/log/supervisor/rag.log
priority=50
startsecs=10

; ── nginx reverse proxy (exposes port 7860) ──────────────────
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log
priority=60
